# ansible/roles/llamacpp/templates/llama_cpp.service.j2
[Unit]
Description=llama.cpp Server (OpenAI-compatible API)
After=network-online.target

[Service]
ExecStart={{ llamacpp_binary_path }}
User={{ llamacpp_user }}
Group={{ llamacpp_group }}
Restart=always
RestartSec=3

Environment="LLAMA_ARG_HOST={{ llamacpp_host }}"
Environment="LLAMA_ARG_PORT={{ llamacpp_port }}"
Environment="LLAMA_ARG_N_GPU_LAYERS={{ llamacpp_n_gpu_layers }}"
Environment="LLAMA_ARG_CTX_SIZE={{ llamacpp_ctx_size }}"
Environment="LLAMA_ARG_N_PARALLEL={{ llamacpp_parallel }}"
Environment="LLAMA_ARG_THREADS={{ llamacpp_threads }}"
{% if llamacpp_model %}
Environment="LLAMA_ARG_MODEL={{ llamacpp_model }}"
{% elif llamacpp_model_hf_repo %}
Environment="LLAMA_ARG_HF_REPO={{ llamacpp_model_hf_repo }}"
{% endif %}
{% if llamacpp_mmproj %}
Environment="LLAMA_ARG_MMPROJ={{ llamacpp_mmproj }}"
{% endif %}

# Allow VRAM overflow to fall back to system RAM instead of crashing
Environment="GGML_CUDA_ENABLE_UNIFIED_MEMORY=1"

Environment="PATH=/usr/local/cuda/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
